# config/systemd/user/dunst.service.d/override.conf
[Unit]
PartOf=
PartOf=i3-session.target

[Service]
Type=
Type=service
BusName=

[Install]
WantedBy=i3-session.target
